# system/etc/profile.d/ratflow-aliases.sh
alias rfcd='cd ~/.config/ratflow/profiles/current/config.d/; pwd; ls -1'
alias rflogout='i3-msg exit'

